# deployment/wc26api.service
[Unit]
Description=WC26 Travel Intel FastAPI
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/home/ec2-user/wc26-travel-intel/backend
Environment="PATH=/home/ec2-user/venv/bin"
ExecStart=/home/ec2-user/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target